# gpu_monitor.service
[Unit]
Description=NVIDIA GPU and LM Studio Token Telemetry Monitor Daemon
After=network.target postgresql.service

[Service]
Type=simple
User=stoneshi
WorkingDirectory=/home/stoneshi/data/src/stone/gpu-monitor
Environment="PATH=/home/stoneshi/.lmstudio/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/usr/bin/python3 /home/stoneshi/data/src/stone/gpu-monitor/gpu_llm_monitor.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
